# deploy/timeweb/yearbook-deploy.service
# Pull-деплой yearbook-v2: oneshot, запускается таймером yearbook-deploy.timer.
# Тянет main с первого доступного источника (primary gitflic → fallback origin/
# github), собирает релиз, переключает симлинк, health-check + откат.
# Ставится в /etc/systemd/system/yearbook-deploy.service
[Unit]
Description=yearbook-v2 pull deploy (gitflic->github main, build, switch, health-check)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=yearbook
Group=yearbook
WorkingDirectory=/srv/yearbook
# Перед деплоем обновляем рабочее дерево репо до свежего main с первого доступного
# источника (gitflic→github) — иначе сам deploy.sh останется на старом коммите.
# Вынесено в prepare-repo.sh (отказоустойчивый выбор источника, общая логика с
# deploy.sh). ExecStartPre завершается ДО ExecStart, скрипт не переписывается на ходу.
ExecStartPre=/srv/yearbook/repo/deploy/timeweb/prepare-repo.sh
ExecStart=/srv/yearbook/repo/deploy/timeweb/deploy.sh
# Сборка тяжёлая — даём время и ограничиваем, чтобы не висло вечно
TimeoutStartSec=1200
Nice=10
# npm пишет в $HOME — у yearbook это /srv/yearbook-v2; зададим явно
Environment=HOME=/srv/yearbook
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
