# deploy/timeweb/yearbook-deploy.service
# Pull-деплой yearbook-v2: oneshot, запускается таймером yearbook-deploy.timer.
# Тянет main с GitHub, собирает релиз, переключает симлинк, health-check + откат.
# Ставится в /etc/systemd/system/yearbook-deploy.service
[Unit]
Description=yearbook-v2 pull deploy (poll GitHub main, build, switch, health-check)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=yearbook
Group=yearbook
WorkingDirectory=/srv/yearbook
# Перед деплоем обновляем рабочее дерево репо-источника до origin/main — иначе
# git fetch обновит только ссылки, а сам deploy.sh останется на старом коммите.
# ExecStartPre завершается ДО ExecStart, поэтому скрипт не переписывается на ходу.
ExecStartPre=/usr/bin/git -C /srv/yearbook/repo fetch --quiet origin main
ExecStartPre=/usr/bin/git -C /srv/yearbook/repo reset --hard --quiet origin/main
ExecStart=/srv/yearbook/repo/deploy/timeweb/deploy.sh
# Сборка тяжёлая — даём время и ограничиваем, чтобы не висло вечно
TimeoutStartSec=1200
Nice=10
# npm пишет в $HOME — у yearbook это /srv/yearbook-v2; зададим явно
Environment=HOME=/srv/yearbook
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
